# packages/base/src/files/etc/cron.d/monitor.cron
* * * * * export POLICY_HOME=${{POLICY_HOME}}; ${{POLICY_HOME}}/bin/monitor >/dev/null 2>&1
